# photobooth.service
[Unit]
Description=PhotoBooth Raspberry Pi
After=network.target graphical.target

[Service]
Type=simple
User=pi
Group=pi

# Répertoire de travail = racine du projet
WorkingDirectory=/home/pi/photobooth

# Activer l'environnement virtuel Python
ExecStart=/home/pi/photobooth/venv/bin/python src/app.py --config config.yaml

# Variables d'environnement nécessaires à Pygame sans bureau
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=SDL_VIDEODRIVER=x11
Environment=PYTHONUNBUFFERED=1

Restart=on-failure
RestartSec=5

# Journalisation
StandardOutput=journal
StandardError=journal
SyslogIdentifier=photobooth

[Install]
WantedBy=graphical.target
